# Linux/tweak.service
[Unit]
Description=Linux System Tweaks

[Service]
ExecStart=/usr/local/bin/tweak
Type=idle

[Install]
WantedBy=multi-user.target
